# agent/wiki-probe-agent.service
[Unit]
Description=Wiki Kele outbound probe agent
# Installed automatically by scripts/install.sh; advanced users may copy this unit manually.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wiki-probe-agent -config /etc/wiki-probe-agent.json
Restart=always
RestartSec=5s
DynamicUser=yes
# Installed by scripts/install.sh for one-line GitHub Release installs.
# ICMP checks still require extra setup; see README. Keep capabilities disabled by default.
# AmbientCapabilities=CAP_NET_RAW
# CapabilityBoundingSet=CAP_NET_RAW
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/etc/wiki-probe-agent.json

[Install]
WantedBy=multi-user.target
